# local.d/scumsnet.stop
#!/bin/sh
cd /root/scumsnet
docker-compose stop

