# local.d/scumsnet.start
#!/bin/sh
cd /root/scumsnet
git pull
docker-compose pull
docker-compose up --no-start
docker-compose start

