# scumsnet/service/deploy/run
#!/bin/sh
sv start scrapyd || exit 1
exec /etc/service/deploy/deploy
