# scumsnet/service/jupyter/run
#!/bin/sh
exec /etc/service/jupyter/jupyter
