# scumsnet/service/scrapydweb/run
#!/bin/sh
sv start scrapyd || exit 1
sv start logparser || exit 1
exec /usr/local/bin/scrapydweb
